# scripts/relay-server.service
[Unit]
Description=RPi Relay Server
After=network.target

[Service]
ExecStart=/home/relayServer
Restart=always

[Install]
WantedBy=multi-user.target
